# 0x1A-application_server/gunicorn.service
[Unit]
Description=starts Gunicorn
After=network.target

[Service]
ExecStart=/startApps.sh

[Install]
WantedBy=multi.user.target
